# xorg.service.in

#
# Minimal Xorg service file - launches Xorg as a service unit
#

# The Xorg launch helper forks, launches Xorg and waits for Xorg to
# accept incoming connections to $DISPLAY, and then signals READY
# to systemd. This guarantees that services that require access to
# $DISPLAY during the session don't start too early.
#
# If you implement a service that requires access to $DISPLAY, your
# service unit file needs to include 'After=xorg.target'.

[Unit]
Description=Xorg server launch helper
Before=xorg.target

[Service]
TTYPath=/dev/tty2
StandardInput=tty-force
EnvironmentFile=-/etc/sysconfig/xorg
Environment=XORG_WATCH_DOG=1
Type=notify
OOMScoreAdjust=-1000
ExecStart=@prefix@/bin/xorg-launch-helper $OPTIONS
